# k6-cron/crontab
*/1 * * * * /usr/local/bin/k6 run /scripts/web_workflow_test.js
